# procursus/procursus/etc/profile.d/cacerts.bootstrap.sh
## git ##
git config --global http.sslCAInfo /opt/procursus/etc/ssl/certs/cacert.pem >/dev/null 2>&1
